# com.kraken.dca.template.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kraken.dca</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/bun</string>
        <string>run</string>
        <string>src/cli.ts</string>
        <string>run</string>
        <!-- Add your config files here -->
        <string>configs/bitcoin-weekly.json</string>
        <string>configs/ethereum-daily.json</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/shivekkhurana/Wip/dkrak</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/shivekkhurana/Wip/dkrak/logs/dca.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/shivekkhurana/Wip/dkrak/logs/dca-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
